# package/alemprator-firstboot/files/etc/init.d/alemprator-firstboot
#!/bin/sh /etc/rc.common

USE_PROCD=1
START=98
EXTRA_COMMANDS="finalize_cleanup monitor_once force_cleanup"
EXTRA_HELP="    finalize_cleanup  Finalize deferred firstboot cleanup\n    monitor_once      Check tracked configs and cleanup if firstboot changed\n    force_cleanup     Immediately cleanup temporary firstboot provisioning"

. /lib/functions.sh

CFG="alemprator_firstboot"
CLEANUP_TOKEN_FILE="/tmp/alemprator-firstboot-cleanup.token"
BASELINE_FILE="/etc/alemprator-firstboot-baseline.md5"
TRACKED_CONFIGS="network wireless dhcp firewall setup"

write_baseline_file() {
	local config
	: > "$BASELINE_FILE"
	for config in $TRACKED_CONFIGS; do
		[ -f "/etc/config/$config" ] || continue
		md5sum "/etc/config/$config" >> "$BASELINE_FILE"
	done
	[ -s "$BASELINE_FILE" ] || rm -f "$BASELINE_FILE"
}

ensure_baseline_file() {
	[ -f "$BASELINE_FILE" ] && return 0
	write_baseline_file
	logger -t alemprator-firstboot "baseline was missing; bootstrap snapshot created"
	return 1
}

baseline_changed() {
	local line
	local digest
	local path
	local current

	[ -f "$BASELINE_FILE" ] || return 1

	while IFS= read -r line; do
		digest="${line%% *}"
		path="${line#*  }"
		[ -n "$path" ] || continue
		[ -f "$path" ] || return 0
		current="$(md5sum "$path" 2>/dev/null | awk '{print $1}')"
		[ -n "$current" ] || return 0
		[ "$current" = "$digest" ] || return 0
	done < "$BASELINE_FILE"

	return 1
}

set_cleanup_pending() {
	uci -q set ${CFG}.main.auto_cleanup_pending='1'
	uci -q commit ${CFG}
}

clear_cleanup_pending() {
	uci -q set ${CFG}.main.auto_cleanup_pending='0'
	uci -q commit ${CFG}
}

load_firstboot_config() {
	config_load "$CFG"
	config_get enabled main enabled 1
	config_get configured_once main configured_once 0
	config_get auto_cleanup_armed main auto_cleanup_armed 0
	config_get auto_cleanup_delay main auto_cleanup_delay 15
	config_get auto_cleanup_pending main auto_cleanup_pending 0
	config_get configured_flag main configured_flag "/etc/configured"
	config_get pending_flag main pending_flag "/etc/alemprator-firstboot-pending"
	config_get network_section main network_section "br_setup"
	config_get dhcp_section main dhcp_section "br_setup"
	config_get wireless_section main wireless_section "alemprator_firstboot"
	config_get firewall_section main firewall_section "br_setup"
}

schedule_deferred_cleanup() {
	local delay
	local token

	delay="$auto_cleanup_delay"
	case "$delay" in
		''|*[!0-9]*) delay=15 ;;
		0) delay=15 ;;
	esac

	token="$(date +%s%N 2>/dev/null || date +%s)$$"
	printf '%s\n' "$token" > "$CLEANUP_TOKEN_FILE"
	uci -q set ${CFG}.main.auto_cleanup_pending='1'
	uci -q commit ${CFG}

	(
		sleep "$delay"
		[ -f "$CLEANUP_TOKEN_FILE" ] || exit 0
		[ "$(cat "$CLEANUP_TOKEN_FILE" 2>/dev/null)" = "$token" ] || exit 0
		/etc/init.d/alemprator-firstboot finalize_cleanup >/dev/null 2>&1 || true
	) &
}

cleanup_firstboot_runtime() {
	local changed_network=0
	local changed_dhcp=0
	local changed_wireless=0
	local changed_firewall=0

	uci -q set ${CFG}.main.enabled='0'
	uci -q set ${CFG}.main.configured_once='1'
	uci -q set ${CFG}.main.auto_cleanup_armed='0'
	uci -q set ${CFG}.main.auto_cleanup_pending='0'
	uci -q commit ${CFG}

	if uci -q get network.${network_section} >/dev/null; then
		uci -q delete network.${network_section}
		changed_network=1
	fi

	if uci -q get dhcp.${dhcp_section} >/dev/null; then
		uci -q delete dhcp.${dhcp_section}
		changed_dhcp=1
	fi

	if uci -q get wireless.${wireless_section} >/dev/null; then
		uci -q delete wireless.${wireless_section}
		changed_wireless=1
	fi

	if uci -q get firewall.${firewall_section} >/dev/null; then
		uci -q delete firewall.${firewall_section}
		changed_firewall=1
	fi

	[ "$changed_network" = 1 ] && uci -q commit network
	[ "$changed_dhcp" = 1 ] && uci -q commit dhcp
	[ "$changed_wireless" = 1 ] && uci -q commit wireless
	[ "$changed_firewall" = 1 ] && uci -q commit firewall

	touch "$configured_flag"
	rm -f "$pending_flag"
	rm -f "$CLEANUP_TOKEN_FILE"
	rm -f "$BASELINE_FILE"

	[ "$changed_network" = 1 ] && /etc/init.d/network reload >/dev/null 2>&1 || true
	if [ "$changed_dhcp" = 1 ]; then
		/etc/init.d/dnsmasq restart >/dev/null 2>&1 || true
		if ! pgrep -x dnsmasq >/dev/null 2>&1; then
			local dnsmasq_cfg
			dnsmasq_cfg="$(uci -q get dhcp.@dnsmasq[0].'.name' 2>/dev/null || echo dnsmasq)"
			/usr/sbin/dnsmasq -C "/var/etc/dnsmasq.conf.${dnsmasq_cfg}" \
				-k -x /var/run/dnsmasq/dnsmasq.pid 2>/dev/null || true
		fi
	fi
	[ "$changed_firewall" = 1 ] && /etc/init.d/firewall restart >/dev/null 2>&1 || true
	[ "$changed_wireless" = 1 ] && wifi reload >/dev/null 2>&1 || true
}

monitor_once() {
	load_firstboot_config

	if [ "$enabled" != "1" ] || [ "$auto_cleanup_armed" != "1" ]; then
		return 0
	fi

	ensure_baseline_file >/dev/null 2>&1 || return 0

	if [ "$(uci -q get setup.default.initial_setup_complete)" = "1" ]; then
		cleanup_firstboot_runtime
		return 0
	fi

	if baseline_changed; then
		[ "$auto_cleanup_pending" = "1" ] || schedule_deferred_cleanup
	else
		[ "$auto_cleanup_pending" = "1" ] && clear_cleanup_pending
	fi
}

finalize_cleanup() {
	load_firstboot_config

	if [ "$enabled" = "1" ] && [ "$auto_cleanup_armed" = "1" ] && [ "$auto_cleanup_pending" = "1" ]; then
		cleanup_firstboot_runtime
	fi
}

force_cleanup() {
	load_firstboot_config

	if [ "$enabled" = "1" ] && [ "$auto_cleanup_armed" = "1" ]; then
		cleanup_firstboot_runtime
	fi
}

sync_firstboot_flags() {
	local complete=0
	local changed=0

	load_firstboot_config

	[ "$(uci -q get setup.default.initial_setup_complete)" = "1" ] && complete=1

	if [ "$enabled" = "0" ]; then
		if [ "$(uci -q get ${CFG}.main.configured_once)" != "1" ]; then
			uci -q set ${CFG}.main.configured_once='1'
			changed=1
		fi

		if [ "$auto_cleanup_armed" != "0" ]; then
			uci -q set ${CFG}.main.auto_cleanup_armed='0'
			changed=1
		fi

		if [ "$auto_cleanup_pending" != "0" ]; then
			uci -q set ${CFG}.main.auto_cleanup_pending='0'
			changed=1
		fi

		[ "$changed" -eq 1 ] && uci -q commit ${CFG}
		touch "$configured_flag"
		rm -f "$pending_flag"
	elif [ "$complete" = "1" ]; then
		cleanup_firstboot_runtime
	else
		touch "$pending_flag"
	fi
}

provision_setup_ap() {
	# Create the temporary Setup AP wireless iface + network + dhcp + firewall
	# if they are missing while firstboot is still pending.  This makes the
	# Setup AP resilient: it is re-created on every boot (not only on the very
	# first uci-defaults run) so a user who reboots before finishing the wizard
	# can always reconnect to the setup network.
	local setup_ip setup_mask radio ssid
	local changed_network=0 changed_dhcp=0 changed_wireless=0 changed_firewall=0

	[ "$enabled" = "1" ] || return 0
	[ -f "$pending_flag" ] || return 0

	setup_ip="$(uci -q get ${CFG}.main.setup_ipaddr 2>/dev/null || echo '192.168.8.1')"
	setup_mask="$(uci -q get ${CFG}.main.setup_netmask 2>/dev/null || echo '255.255.255.0')"
	radio="$(uci -q get ${CFG}.main.radio 2>/dev/null || echo 'radio0')"
	ssid="$(uci -q get ${CFG}.main.ssid 2>/dev/null || echo 'ALemprator-Setup')"

	# Setup AP network interface: a direct L3 interface bound to the wifi-iface
	# (no bridge device, no explicit 'device' — the wireless binding with
	# network='${network_section}' creates the L3 link on phy0-ap0). Setting an
	# explicit device name forced netifd to look for a non-existent netdev.
	if ! uci -q get network.${network_section} >/dev/null 2>&1; then
		uci -q set network.${network_section}=interface
		changed_network=1
	fi
	if [ -n "$(uci -q get network.${network_section}.device)" ]; then
		uci -q delete network.${network_section}.device
		changed_network=1
	fi
	if [ "$(uci -q get network.${network_section}.proto)" != "static" ]; then
		uci -q set network.${network_section}.proto='static'
		changed_network=1
	fi
	if [ "$(uci -q get network.${network_section}.ipaddr)" != "${setup_ip}" ]; then
		uci -q set network.${network_section}.ipaddr="${setup_ip}"
		changed_network=1
	fi
	if [ "$(uci -q get network.${network_section}.netmask)" != "${setup_mask}" ]; then
		uci -q set network.${network_section}.netmask="${setup_mask}"
		changed_network=1
	fi
	[ "$changed_network" = 1 ] && uci -q commit network

	# DHCP section
	if ! uci -q get dhcp.${dhcp_section} >/dev/null 2>&1; then
		uci -q batch <<EOF >/dev/null 2>&1
			set dhcp.${dhcp_section}=dhcp
			set dhcp.${dhcp_section}.interface='${network_section}'
			set dhcp.${dhcp_section}.start='100'
			set dhcp.${dhcp_section}.limit='50'
			set dhcp.${dhcp_section}.leasetime='1h'
			delete dhcp.${dhcp_section}.ignore
			commit dhcp
EOF
		changed_dhcp=1
	fi

	# Firewall section
	if ! uci -q get firewall.${firewall_section} >/dev/null 2>&1; then
		uci -q batch <<EOF >/dev/null 2>&1
			set firewall.${firewall_section}=zone
			set firewall.${firewall_section}.name='${firewall_section}'
			set firewall.${firewall_section}.network='${network_section}'
			set firewall.${firewall_section}.input='ACCEPT'
			set firewall.${firewall_section}.output='ACCEPT'
			set firewall.${firewall_section}.forward='REJECT'
			commit firewall
EOF
		changed_firewall=1
	fi

	# Wireless Setup AP iface (direct L3 binding, no bridge)
	if ! uci -q get wireless.${wireless_section} >/dev/null 2>&1; then
		uci -q batch <<EOF >/dev/null 2>&1
			set wireless.${wireless_section}=wifi-iface
			set wireless.${wireless_section}.device='${radio}'
			set wireless.${wireless_section}.mode='ap'
			set wireless.${wireless_section}.ssid='${ssid}'
			set wireless.${wireless_section}.encryption='none'
			set wireless.${wireless_section}.network='${network_section}'
			set wireless.${wireless_section}.disabled='0'
			set wireless.${wireless_section}.isolate='1'
			commit wireless
EOF
		changed_wireless=1
	fi
	# Ensure the wifi iface is bound to the radio and the setup network
	if [ "$(uci -q get wireless.${wireless_section}.device)" != "${radio}" ]; then
		uci -q set wireless.${wireless_section}.device="${radio}"
		changed_wireless=1
	fi
	uci -q set wireless.${wireless_section}.network="${network_section}"
	uci -q commit wireless

	[ "$changed_network" = 1 ] && /etc/init.d/network restart >/dev/null 2>&1 || true
	[ "$changed_firewall" = 1 ] && /etc/init.d/firewall restart >/dev/null 2>&1 || true
	[ "$changed_wireless" = 1 ] && wifi reload >/dev/null 2>&1 || true
	# DHCP is mandatory on the setup network so clients can obtain an IP and
	# reach the captive portal. Always (re)start dnsmasq when provisioning.
	/etc/init.d/dnsmasq restart >/dev/null 2>&1 || true
	if ! pgrep -x dnsmasq >/dev/null 2>&1; then
		# init script restart may fail silently (e.g. udhcpc hook); fall back
		# to starting dnsmasq directly from the generated config file. Do not
		# trust dhcp.@dnsmasq[0].'.name' (it can drift from the generated file
		# name) — locate the actual generated config under /var/etc instead.
		local dnsmasq_conf
		dnsmasq_conf="$(ls -1 /var/etc/dnsmasq.conf.* 2>/dev/null | head -1)"
		[ -f "$dnsmasq_conf" ] || dnsmasq_conf="/var/etc/dnsmasq.conf"
		mkdir -p /var/run/dnsmasq
		/usr/sbin/dnsmasq -C "$dnsmasq_conf" \
			-k -x /var/run/dnsmasq/dnsmasq.pid 2>/dev/null || true
	fi

	[ "$changed_network" = 1 ] || [ "$changed_dhcp" = 1 ] || [ "$changed_wireless" = 1 ] || [ "$changed_firewall" = 1 ]
}

boot() {
	load_firstboot_config

	# Factory-reset safety net: if the device is not yet configured OR the
	# setup section names are stale (e.g. uci-defaults did not run, or an old
	# config carried over after firstboot -y wiped the overlay), bootstrap the
	# firstboot flags and force the canonical br_setup section names so the
	# Setup AP is always created on a clean device.
	if [ "$enabled" != "1" ] || [ "$configured_once" != "1" ] || [ ! -f "$configured_flag" ] || \
	   [ "$network_section" != "br_setup" ] || [ "$dhcp_section" != "br_setup" ] || [ "$firewall_section" != "br_setup" ]; then
		uci -q set ${CFG}.main.enabled='1'
		uci -q set ${CFG}.main.configured_once='0'
		uci -q set ${CFG}.main.auto_cleanup_armed='1'
		uci -q set ${CFG}.main.network_section='br_setup'
		uci -q set ${CFG}.main.dhcp_section='br_setup'
		uci -q set ${CFG}.main.firewall_section='br_setup'
		uci -q commit ${CFG}
		touch "$pending_flag"
		load_firstboot_config
	fi

	if [ "$enabled" = "1" ] && [ "$auto_cleanup_armed" = "1" ] && [ "$auto_cleanup_pending" = "1" ]; then
		cleanup_firstboot_runtime
	fi

	sync_firstboot_flags
	load_firstboot_config

	if [ "$enabled" = "1" ] && [ "$auto_cleanup_armed" = "1" ] && [ -f "$pending_flag" ]; then
		provision_setup_ap >/dev/null 2>&1 || true
		ensure_baseline_file >/dev/null 2>&1 || true
		setup_captive_portal
	fi
}

setup_captive_portal() {
	# IP شبكة الإعداد المؤقتة
	local setup_ip
	setup_ip="$(uci -q get network.${network_section}.ipaddr 2>/dev/null || echo '192.168.8.1')"

	# DNS: أعد جميع الاستعلامات لـ IP الراوتر (يُفعّل Captive Portal)
	uci -q set dhcp.alemprator_captive=dnsmasq_option 2>/dev/null || true
	uci -q batch <<EOF >/dev/null 2>&1
		set dhcp.@dnsmasq[0].address='/#/${setup_ip}'
		commit dhcp
EOF
	/etc/init.d/dnsmasq restart >/dev/null 2>&1 || true
	# Fallback: start directly if jail fails
	if ! pgrep -x dnsmasq >/dev/null 2>&1; then
		local dnsmasq_cfg
		dnsmasq_cfg="$(uci -q get dhcp.@dnsmasq[0].'.name' 2>/dev/null || echo dnsmasq)"
		/usr/sbin/dnsmasq -C "/var/etc/dnsmasq.conf.${dnsmasq_cfg}" \
			-k -x /var/run/dnsmasq/dnsmasq.pid 2>/dev/null || true
	fi

	# nftables: أعد HTTP (80) من العملاء للراوتر
	nft add chain inet fw4 alemprator_captive 2>/dev/null || true
	nft flush chain inet fw4 alemprator_captive 2>/dev/null || true
	nft add rule inet fw4 dstnat \
		iifname "${network_section}" \
		ip daddr != "${setup_ip}" \
		tcp dport 80 \
		dnat to "${setup_ip}:80" 2>/dev/null || true

	logger -t alemprator-firstboot "Captive portal activated on ${setup_ip}"
}

remove_captive_portal() {
	# إزالة قاعدة DNS الخاصة بالـ Captive Portal
	uci -q delete dhcp.@dnsmasq[0].address 2>/dev/null || true
	uci -q commit dhcp 2>/dev/null || true
	/etc/init.d/dnsmasq restart >/dev/null 2>&1 || true
	# Fallback: start directly if jail fails
	if ! pgrep -x dnsmasq >/dev/null 2>&1; then
		local dnsmasq_cfg
		dnsmasq_cfg="$(uci -q get dhcp.@dnsmasq[0].'.name' 2>/dev/null || echo dnsmasq)"
		/usr/sbin/dnsmasq -C "/var/etc/dnsmasq.conf.${dnsmasq_cfg}" \
			-k -x /var/run/dnsmasq/dnsmasq.pid 2>/dev/null || true
	fi
	logger -t alemprator-firstboot "Captive portal removed"
}

start_service() {
	load_firstboot_config

	if [ "$enabled" = "1" ] && [ "$auto_cleanup_armed" = "1" ] && [ "$auto_cleanup_pending" = "1" ]; then
		cleanup_firstboot_runtime
	fi

	sync_firstboot_flags
	load_firstboot_config

	if [ "$enabled" = "1" ] && [ "$auto_cleanup_armed" = "1" ] && [ -f "$pending_flag" ]; then
		provision_setup_ap >/dev/null 2>&1 || true
		ensure_baseline_file >/dev/null 2>&1 || true
		setup_captive_portal
	fi

	if [ "$enabled" = "1" ] && [ "$auto_cleanup_armed" = "1" ]; then
		procd_open_instance
		procd_set_param command /bin/sh -c 'while true; do /etc/init.d/alemprator-firstboot monitor_once >/dev/null 2>&1 || true; sleep 5; done'
		procd_close_instance
	fi
}


reload_service() {
	load_firstboot_config

	if [ "$enabled" = "1" ] && [ "$auto_cleanup_armed" = "1" ]; then
		ensure_baseline_file >/dev/null 2>&1 || true
		schedule_deferred_cleanup
	fi

	sync_firstboot_flags
}

service_triggers() {
	procd_add_reload_trigger "alemprator_firstboot" "network" "wireless" "dhcp" "firewall" "system" "setup"
}