# ota-server/public/firmware/_openwrt-ramips-mt7621-kt_km14-102h-squashfs-sysupgrade-24.10.4-km14-r41.bin.extracted/sysupgrade-kt_km14-102h/squashfs-root/etc/hotplug.d/iface/95-alemprator-ota
#!/bin/sh

[ "$ACTION" = "ifup" ] || exit 0
case "$INTERFACE" in
	wan|wwan|lan) ;;
	*) exit 0 ;;
esac

/usr/libexec/alemprator-ota/run-once --check-only >/dev/null 2>&1 &
exit 0
