# inspect_ipk/usr/libexec/hotspot-openwrt/init.d/hotspot-openwrt
#!/bin/sh /etc/rc.common

START=91
USE_PROCD=1

start_service() {
	[ "$(uci -q get hotspot_openwrt.main.enabled)" = '1' ] || return 0
	/usr/libexec/hotspot-openwrt/apply --start
}

stop_service() {
	if [ -s /var/run/hotspot-openwrt-presence-watchdog.pid ]; then
		kill "$(cat /var/run/hotspot-openwrt-presence-watchdog.pid 2>/dev/null)" 2>/dev/null || true
		rm -f /var/run/hotspot-openwrt-presence-watchdog.pid
	fi
	if [ -s /var/run/hotspot-openwrt-mac-cookie-watchdog.pid ]; then
		kill "$(cat /var/run/hotspot-openwrt-mac-cookie-watchdog.pid 2>/dev/null)" 2>/dev/null || true
		rm -f /var/run/hotspot-openwrt-mac-cookie-watchdog.pid
	fi
	/etc/init.d/chilli stop 2>/dev/null || true
}