# ota-server/km15-103h-aftar-upgrade/extracted_old_squashfs/usr/lib/opkg/info/mcr-busybox.list
/bin/dmesg
/bin/login
/bin/rm
/sbin/ip
/sbin/lsmod
/bin/sh
/sbin/vconfig
/usr/bin/test
/sbin/depmod
/bin/zcat
/bin/fgrep
/bin/ash
/sbin/halt
/sbin/reboot
/bin/true
/sbin/route
/usr/bin/time
/sbin/klogd
/usr/bin/traceroute
/bin/hostname
/bin/vi
/bin/pwd
/sbin/modprobe
/usr/bin/renice
/etc/init.d/cron
/bin/mount
/bin/sync
/sbin/rmmod
/sbin/insmod
/bin/dnsdomainname
/bin/echo
/usr/sbin/crond
/usr/bin/head
/bin/gunzip
/bin/netstat
/bin/busybox
/bin/grep
/bin/sleep
/bin/kill
/bin/mktemp
/bin/mknod
/usr/sbin/massp
/usr/bin/mcrfet
/usr/bin/hexdump
/etc/init.d/sysntpd
/bin/date
/usr/bin/wc
/usr/bin/[
/bin/umount
/usr/bin/tail
/bin/mpstat
/usr/sbin/ntpd-hotplug
/bin/cp
/sbin/poweroff
/usr/bin/[[
/bin/dd
/bin/lock
/bin/df
/usr/bin/killall
/bin/mkdir
/sbin/ifconfig
/bin/cat
/usr/bin/lzcat
/bin/touch
/bin/tar
/usr/bin/taskset
/usr/bin/top
/usr/bin/expr
/sbin/sysctl
/sbin/syslogd
/bin/sed
/bin/chmod
/usr/bin/awk
/bin/uname
/usr/bin/which
/bin/ln
/bin/ls
/usr/bin/cut
/bin/mv
/usr/bin/crontab
/usr/bin/tty
/bin/ping
/usr/bin/uptime
/bin/fsync
/bin/egrep
/usr/bin/readlink
/usr/sbin/udhcpc
/usr/sbin/udhcpd
/bin/ps
/bin/gzip
/usr/bin/dumpleases
/usr/bin/unlzma
