# package/luci-app-alemprator-ota/files/etc/uci-defaults/95-alemprator-ota
#!/bin/sh

model_file="/etc/model"
version_file="/etc/alemprator/firmware-version"
identity_file="/etc/alemprator/model-identities"
board_name="$(ubus call system board 2>/dev/null | jsonfilter -e '@.board_name')"
model_key=""
firmware_version=""

read_identity_field() {
	local field field_index
	field="$1"

	[ -n "$board_name" ] || return 1
	[ -s "$identity_file" ] || return 1

	case "$field" in
		model_key) field_index=2 ;;
		firmware_version) field_index=3 ;;
		*) return 1 ;;
	esac

	awk -F '|' -v board="$board_name" -v field_index="$field_index" '
		$1 == board && $0 !~ /^#/ {
			print $field_index
			found = 1
			exit
		}
		END {
			if (found)
				exit 0
			exit 1
		}
	' "$identity_file"
}

model_key="$(read_identity_field model_key 2>/dev/null || true)"
firmware_version="$(read_identity_field firmware_version 2>/dev/null || true)"

if [ -n "$model_key" ]; then
	current_model="$(cat "$model_file" 2>/dev/null)"
	[ "$current_model" = "$model_key" ] || printf '%s\n' "$model_key" > "$model_file"
elif [ ! -s "$model_file" ]; then
	[ -n "$board_name" ] || board_name="unknown"
	printf '%s\n' "$board_name" > "$model_file"
fi

if [ -n "$firmware_version" ]; then
	current_version="$(cat "$version_file" 2>/dev/null)"
	[ "$current_version" = "$firmware_version" ] || printf '%s\n' "$firmware_version" > "$version_file"
fi

uci -q batch <<EOF >/dev/null
	set alemprator_ota.main=ota
	set alemprator_ota.main.enabled='1'
	set alemprator_ota.main.server_url='https://ota.kartnet.org'
	set alemprator_ota.main.register_path='/api/register'
	set alemprator_ota.main.update_path='/api/update'
	set alemprator_ota.main.heartbeat_path='/api/heartbeat'
	set alemprator_ota.main.check_interval='21600'
	set alemprator_ota.main.random_delay_max='3600'
	set alemprator_ota.main.auto_upgrade='1'
	set alemprator_ota.main.require_setup_complete='0'
	set alemprator_ota.main.keep_config='1'
	set alemprator_ota.main.allow_force='1'
	set alemprator_ota.main.window_start='2'
	set alemprator_ota.main.window_end='6'
	set alemprator_ota.main.retry_base='900'
	set alemprator_ota.main.retry_max='21600'
	set alemprator_ota.main.token_salt='CHANGE_ME_UNIQUE_PER_BRAND'
	set alemprator_ota.main.hmac_secret=''
	set alemprator_ota.main.connect_timeout='20'
	set alemprator_ota.main.model_file='/etc/model'
	set alemprator_ota.main.model_identity_file='/etc/alemprator/model-identities'
	set alemprator_ota.main.token_file='/etc/alemprator/device.token'
	set alemprator_ota.main.version_file='/etc/alemprator/firmware-version'
	set alemprator_ota.main.state_file='/tmp/alemprator-ota/state.env'
	commit alemprator_ota
EOF

/etc/init.d/alemprator-ota enable >/dev/null 2>&1 || true
/etc/init.d/alemprator-ota restart >/dev/null 2>&1 || true

(
	sleep 30
	/usr/libexec/alemprator-ota/run-once >/dev/null 2>&1 || true
) &

exit 0
