# package/luci-app-alemprator-ota/files/etc/uci-defaults/95-alemprator-ota
#!/bin/sh

model_file="/etc/model"
board_name="$(ubus call system board 2>/dev/null | jsonfilter -e '@.board_name')"

if [ ! -s "$model_file" ]; then
	[ -n "$board_name" ] || board_name="unknown"
	printf '%s\n' "$board_name" > "$model_file"
fi

if ! uci -q get alemprator_ota.main >/dev/null; then
	uci -q batch <<EOF >/dev/null
		set alemprator_ota.main=ota
		set alemprator_ota.main.enabled='1'
		set alemprator_ota.main.server_url='https://updates.example.com'
		set alemprator_ota.main.register_path='/api/register'
		set alemprator_ota.main.update_path='/api/update'
		set alemprator_ota.main.heartbeat_path='/api/heartbeat'
		set alemprator_ota.main.check_interval='21600'
		set alemprator_ota.main.random_delay_max='3600'
		set alemprator_ota.main.auto_upgrade='1'
		set alemprator_ota.main.keep_config='1'
		set alemprator_ota.main.allow_force='1'
		set alemprator_ota.main.window_start='2'
		set alemprator_ota.main.window_end='6'
		set alemprator_ota.main.retry_base='900'
		set alemprator_ota.main.retry_max='21600'
		set alemprator_ota.main.token_salt='CHANGE_ME_UNIQUE_PER_BRAND'
		set alemprator_ota.main.hmac_secret=''
		set alemprator_ota.main.connect_timeout='20'
		set alemprator_ota.main.model_file='/etc/model'
		set alemprator_ota.main.token_file='/etc/alemprator/device.token'
		set alemprator_ota.main.version_file='/etc/alemprator/firmware-version'
		set alemprator_ota.main.state_file='/tmp/alemprator-ota/state.env'
		commit alemprator_ota
EOF
fi

/etc/init.d/alemprator-ota enable >/dev/null 2>&1 || true
/etc/init.d/alemprator-ota restart >/dev/null 2>&1 || true

exit 0
