# ota-server/public/firmware/_openwrt-ramips-mt7621-kt_km14-102h-squashfs-sysupgrade-24.10.4-km14-r41.bin.extracted/sysupgrade-kt_km14-102h/squashfs-root/usr/lib/opkg/info/alemprator-firstboot.list
/etc/init.d/alemprator-firstboot
/etc/uci-defaults/95-alemprator-firstboot
/usr/share/ucitrack/alemprator-firstboot-dhcp.json
/etc/config/alemprator_firstboot
/usr/share/ucitrack/alemprator-firstboot-wireless.json
/usr/share/ucitrack/alemprator-firstboot-firewall.json
/usr/share/ucitrack/alemprator-firstboot-network.json
/usr/share/ucitrack/alemprator-firstboot-setup.json
/etc/rc.d/S98alemprator-firstboot
