# ota-server/public/firmware/_openwrt-ramips-mt7621-kt_km14-102h-squashfs-sysupgrade-24.10.4-km14-r41.bin.extracted/sysupgrade-kt_km14-102h/squashfs-root/usr/lib/opkg/info/bandix-plus.list
/etc/config/bandix-plus
/usr/bin/bandix-plus
/lib/upgrade/keep.d/bandix-plus
/etc/init.d/bandix-plus
