# ota-server/public/firmware/_openwrt-ramips-mt7621-kt_km14-102h-squashfs-sysupgrade-24.10.4-km14-r41.bin.extracted/sysupgrade-kt_km14-102h/squashfs-root/usr/lib/opkg/info/coova-chilli.list
/usr/sbin/chilli_query
/etc/chilli/www/chilliform.chi.sh
/etc/chilli/www/ChilliLibrary.js
/usr/sbin/chilli
/etc/chilli/www/chilliform.chi
/etc/init.d/chilli
/etc/chilli/gui-config-default.ini
/etc/chilli/www/coova.html
/usr/sbin/chilli_response
/etc/chilli/www/coova.png
/etc/chilli/www/json_html.tmpl
/etc/chilli/www/chillijs.chi.sh
/etc/config/chilli
/etc/chilli/functions
/etc/chilli/www/wait.gif
/etc/chilli/wpad.dat
/etc/chilli/defaults
/usr/lib/libchilli.so.0.0.0
/usr/sbin/chilli_opt
/usr/lib/libbstring.so.0
/etc/chilli/www/chillijs.chi
/usr/lib/libbstring.so.0.0.0
/usr/sbin/chilli_radconfig
/etc/chilli/wwwsh
/etc/chilli/www/openid.gif
/etc/chilli/newmulti.sh
/etc/chilli/www/config.sh
/usr/lib/libchilli.so.0
/etc/chilli/www/coova.jpg
/etc/chilli/www/chilliController.js
/etc/chilli.conf
