# ota-server/public/firmware/_openwrt-ramips-mt7621-kt_km14-102h-squashfs-sysupgrade-24.10.4-km14-r41.bin.extracted/sysupgrade-kt_km14-102h/squashfs-root/usr/lib/opkg/info/librespeed-go.list
/usr/bin/librespeed-go
/lib/upgrade/keep.d/librespeed-go
/etc/init.d/librespeed-go
/etc/librespeed-go/telemetry_postgresql.sql
/etc/librespeed-go/telemetry_mysql.sql
/etc/config/librespeed-go
