# ota-server/public/firmware/_openwrt-ramips-mt7621-kt_km14-102h-squashfs-sysupgrade-24.10.4-km14-r41.bin.extracted/sysupgrade-kt_km14-102h/squashfs-root/usr/lib/opkg/info/lldpd.list
/etc/config/lldpd
/usr/sbin/lldpd
/usr/sbin/lldpcli
/usr/lib/liblldpctl.so.4.9.1
/usr/lib/liblldpctl.so
/etc/init.d/lldpd
/usr/lib/liblldpctl.so.4
/usr/sbin/lldpctl
