# ota-server/public/firmware/_openwrt-ramips-mt7621-kt_km14-102h-squashfs-sysupgrade-24.10.4-km14-r41.bin.extracted/sysupgrade-kt_km14-102h/squashfs-root/usr/lib/opkg/info/luci-app-netspeedtest.list
/usr/share/netspeedtest/openspeedtest/assets/images/icons/browserconfig.xml
/usr/share/netspeedtest/openspeedtest/assets/fonts/roboto-v30-latin-regular.woff2
/usr/share/netspeedtest/openspeedtest/assets/fonts/roboto-v30-latin-500.woff2
/usr/share/netspeedtest/openspeedtest/assets/images/icons/site.webmanifest
/etc/uci-defaults/99_netspeedtest
/usr/share/netspeedtest/openspeedtest/assets/css/app.css
/usr/share/netspeedtest/openspeedtest/assets/images/app.svg
/usr/share/netspeedtest/openspeedtest/cgi-bin/upload
/usr/share/netspeedtest/openspeedtest/assets/css/darkmode.css
/usr/share/netspeedtest/openspeedtest/License.md
/usr/share/netspeedtest/openspeedtest/assets/images/icons/launcher-icon-1x.png
/usr/share/netspeedtest/openspeedtest/assets/images/icons/android-chrome-192x192.png
/usr/share/netspeedtest/openspeedtest/assets/fonts/roboto-v30-latin-500.woff
/etc/init.d/netspeedtest
/usr/share/netspeedtest/openspeedtest/assets/images/icons/safari-pinned-tab.svg
/www/luci-static/resources/view/netspeedtest/iperf3.js
/usr/share/netspeedtest/openspeedtest/assets/fonts/roboto-v30-latin-regular.woff
/usr/share/netspeedtest/openspeedtest/assets/images/icons/favicon-16x16.png
/usr/share/netspeedtest/openspeedtest/assets/images/icons/launcher-icon-2x.png
/www/luci-static/resources/view/netspeedtest/speedtest.js
/www/luci-static/resources/view/netspeedtest/webspeedtest.js
/usr/share/netspeedtest/openspeedtest/assets/js/darkmode.js
/usr/share/netspeedtest/openspeedtest/index.html
/etc/config/netspeedtest
/usr/share/netspeedtest/openspeedtest/assets/images/icons/mstile-150x150.png
/usr/share/netspeedtest/openspeedtest/assets/images/icons/launcher-icon-3x.png
/usr/share/netspeedtest/openspeedtest/assets/fonts/roboto-v30-latin-500.eot
/usr/share/rpcd/acl.d/luci-app-netspeedtest.json
/usr/share/netspeedtest/openspeedtest/assets/images/icons/favicon.ico
/usr/share/netspeedtest/openspeedtest/assets/js/app-2.5.4.min.js
/usr/share/netspeedtest/openspeedtest/assets/fonts/roboto-v30-latin-500.ttf
/usr/share/netspeedtest/openspeedtest/assets/fonts/roboto-v30-latin-500.svg
/usr/libexec/rpcd/luci.netspeedtest
/usr/share/netspeedtest/openspeedtest/assets/images/icons/launcher-icon-4x.png
/usr/share/netspeedtest/openspeedtest/upload
/usr/share/netspeedtest/openspeedtest/assets/images/icons/android-chrome-512x512.png
/www/luci-static/resources/view/netspeedtest/librespeed.js
/usr/share/netspeedtest/openspeedtest/assets/fonts/roboto-v30-latin-regular.eot
/usr/share/luci/menu.d/luci-app-netspeedtest.json
/usr/share/netspeedtest/openspeedtest/assets/images/icons/favicon-32x32.png
/usr/share/netspeedtest/openspeedtest/assets/fonts/roboto-v30-latin-regular.ttf
/usr/share/netspeedtest/openspeedtest/assets/fonts/roboto-v30-latin-regular.svg
/usr/share/netspeedtest/openspeedtest/assets/images/icons/apple-touch-icon.png
