# ota-server/public/firmware/_openwrt-ramips-mt7621-kt_km14-102h-squashfs-sysupgrade-24.10.4-km14-r41.bin.extracted/sysupgrade-kt_km14-102h/squashfs-root/usr/lib/opkg/info/luci-app-setup.list
/www/luci-static/resources/view/setup/setup.js
/etc/uci-defaults/40_luci-app-setup
/usr/libexec/alemprator-setup/cleanup-hotspot
/etc/uci-defaults/45_luci-app-setup-migrate-wizardvlan
/etc/rc.d/S99setup
/www/captive-portal.html
/etc/uci-defaults/49_luci-app-setup-apply-radio-htmode
/etc/config/setup
/usr/share/luci/menu.d/zzz-luci-app-setup.json
/etc/uci-defaults/46_luci-app-setup-disable-low-ack
/usr/libexec/alemprator-setup/reset-disabled
/etc/uci-defaults/47_luci-app-setup-fix-5g-vlan-ssid
/www/luci-static/resources/view/setup/setup_r93.js
/etc/uci-defaults/48_luci-app-setup-recover-state
/etc/init.d/setup
/etc/uci-defaults/alemprator_c
/usr/share/ucitrack/luci-app-setup.json
/usr/share/rpcd/acl.d/luci-app-setup.json
/www/luci-static/resources/view/setup/landing_r2.js
/www/luci-static/resources/view/setup/landing.js
/usr/libexec/alemprator-setup/wps-disabled
