# ota-server/public/firmware/_openwrt-ramips-mt7621-kt_km14-102h-squashfs-sysupgrade-24.10.4-km14-r41.bin.extracted/sysupgrade-kt_km14-102h/squashfs-root/usr/lib/opkg/info/luci-app-hotspot-openwrt.list
/etc/chilli/condown.sh
/usr/libexec/hotspot-openwrt/mac-cookie-watchdog
/usr/libexec/hotspot-openwrt/apply
/usr/libexec/hotspot-openwrt/logs
/www/cgi-bin/generate_204
/www/hotspot/style.css
/etc/chilli/up.sh
/etc/init.d/hotspot-openwrt
/usr/libexec/hotspot-openwrt/gen-ssl
/www/hotspot/api.json
/usr/libexec/hotspot-openwrt/connectivity-test
/usr/share/rpcd/acl.d/luci-app-hotspot-openwrt.json
/etc/chilli/conup.sh
/usr/libexec/hotspot-openwrt/userman-info
/usr/share/luci/menu.d/luci-app-hotspot-openwrt.json
/www/hotspot/terms.html
/usr/libexec/hotspot-openwrt/kick-client
/etc/hotspot-openwrt/firewall.sh
/usr/libexec/hotspot-openwrt/phase3-rc-gate
/usr/libexec/hotspot-openwrt/import-config
/usr/libexec/hotspot-openwrt/presence-watchdog
/www/cgi-bin/hotspot-speedtest
/www/hotspot/status.html
/etc/chilli/down.sh
/www/cgi-bin/hotspot-logout
/usr/libexec/hotspot-openwrt/status-json
/www/cgi-bin/hotspot-card-info
/www/hotspot/login.html
/etc/uci-defaults/99_hotspot-openwrt-firewall
/www/cgi-bin/hotspot-login
/usr/libexec/hotspot-openwrt/phase2-client-sim
/usr/libexec/hotspot-openwrt/validate
/usr/libexec/hotspot-openwrt/export-config
/usr/libexec/hotspot-openwrt/mac-cookies
/www/luci-static/resources/view/services/hotspot-openwrt.js
/www/hotspot/settings.js
/usr/libexec/hotspot-openwrt/license-check
/etc/config/hotspot_openwrt
/etc/config/hotspot_licensing
/usr/libexec/hotspot-openwrt/portal-upload
/usr/libexec/hotspot-openwrt/phase2-smoke
